# scripts/launchd/com.autonomous-dev.triage.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.autonomous-dev.triage</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd /Users/akaszubski/Dev/autonomous-dev && bash scripts/drain-all.sh --root-cause --limit 6 --timeout 5400</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/akaszubski/Dev/autonomous-dev/logs/scheduled/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/akaszubski/Dev/autonomous-dev/logs/scheduled/launchd.err.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/akaszubski/Dev/autonomous-dev</string>

    <!-- launchd will skip ticks if the prior run is still active -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Inherit user PATH so claude, gh, python3, git are found -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
